# kubeadm-aio/assets/usr/lib/systemd/system/container-up.target
[Unit]
Description=Minimal target for containerized applications
DefaultDependencies=false
AllowIsolate=yes
Requires=systemd-tmpfiles-setup.service systemd-journald.service dbus.service
After=systemd-tmpfiles-setup.service systemd-journald.service dbus.service
