# slackbridge/conf/supervisord.conf
[supervisord]
logfile=/tmp/supervisord.log 
logfile_maxbytes=50MB        
logfile_backups=10           
loglevel=info                
pidfile=/tmp/supervisord.pid 
nodaemon=false              
minfds=1024                 
minprocs=200 

[program:slack-irc-config]
priority=1
command=/bin/bash /slack-irc/slack-irc-config.sh 
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log
autorestart=false
startsecs=0

[program:slack-irc]
priority=30
directory=/slack-irc
command=slack-irc --config /slack-irc/config.json
user=root
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/var/log/slack-irc.log
stderr_logfile=/var/log/slack-irc.err

[unix_http_server]
file=%(here)s/supervisor.sock

[supervisorctl]
serverurl=unix://%(here)s/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface